# ops/mina-binance-listings.service
[Unit]
Description=MINA Binance Futures Listing Watcher
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/MINA_v2
ExecStart=/root/MINA_v2/venv/bin/python -u /root/MINA_v2/signal_bot/binance_listings_watcher.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/root/MINA_v2/signal_bot/binance_listings_watcher.log
StandardError=append:/root/MINA_v2/signal_bot/binance_listings_watcher.log

[Install]
WantedBy=multi-user.target
